# niri/.config/systemd/user/batt.service
[Unit]
Description=Battery monitor script

[Service]
Type=oneshot
ExecStart=%h/.bin/batt
